# overlay/etc/apk/protected_paths.d/lbu.list
+usr/include/nng/
+usr/lib/libnng.so.1.6.0-pre
+usr/lib/libnng.so.1
+usr/lib/libnng.so
+usr/lib/cmake/nng/nng-targets.cmake
+usr/lib/cmake/nng/nng-targets-release.cmake
+usr/lib/cmake/nng/nng-config.cmake
+usr/lib/cmake/nng/nng-config-version.cmake
+usr/bin/nngcat
+usr/share/satdump/
+usr/lib/libsatdump_core.so
+usr/include/satdump/
+usr/bin/satdump
+usr/lib/satdump/
+usr/bin/satdump_sdr_server
+usr/share/satdump/pipelines
+usr/share/applications/satdump.desktop
+root/autotrack.json
+etc/init.d/satdump
